# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:gunicorn]
command=gunicorn "app:create_app('production')" -c /app/gunicorn.conf.py
directory=/app